# .docker/config/s6/fluent-bit/run
#!/command/with-contenv sh
/opt/fluent-bit/bin/fluent-bit -c /etc/fluent-bit/fluentbit.conf
