# .docker/config/s6/queue/run
#!/command/with-contenv sh
php /var/www/artisan queue:work --sleep=3 --tries=3 --backoff=3 --max-time=3600
