# .docker/config/s6/fpm/run
#!/command/with-contenv sh
set -e
# Start PHP-FPM
/usr/sbin/php-fpm8.1 -R --nodaemonize
